# dist/futaba.service
[Unit]
Description=A silly Discord bot for a friend.
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/bin/futaba
Restart=always

[Install]
WantedBy=multi-user.target
